# scripts/cfgsync.service
[Unit]
Description=cfgsync: software configuration sync backend
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=cfgsync
Group=cfgsync
WorkingDirectory=/opt/cfgsync
ExecStart=/opt/cfgsync/cfgsync
Restart=on-failure
RestartSec=5s
EnvironmentFile=/etc/cfgsync/env
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/cfgsync/data /opt/cfgsync/backups
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
